# configuration/files/battery_warning.service
[Unit]
Description=Execute battery warning script

[Service]
ExecStart=/opt/bin/battery_warning.sh

